# sftp-mirror.timer
[Unit]
Description=SFTP Mirror Timer - Start at 14:30 UTC
#Documentation=

[Timer]
# Start the service daily at exactly 14:30 UTC
OnCalendar=*-*-* 14:30:00 UTC
Unit=sftp-mirror.service

# Run immediately after boot if it missed the last scheduled run
Persistent=false

# No randomization - start at exactly 14:30 UTC
RandomizedDelaySec=0

[Install]
WantedBy=timers.target

